# vnstat/vnstat.ini
[program:vnstat]
command=/opt/go/bin/go run /tmp/vnstat.go
environment=GOROOT=/opt/go
directory=/tmp